# quiz-websocket.service
# Quiz Application WebSocket Service
# Systemd service file for running WebSocket server in production
#
# Installation:
# 1. Copy this file to /etc/systemd/system/quiz-websocket.service
# 2. Edit paths to match your installation
# 3. Run: sudo systemctl daemon-reload
# 4. Run: sudo systemctl enable quiz-websocket
# 5. Run: sudo systemctl start quiz-websocket
# 6. Check status: sudo systemctl status quiz-websocket

[Unit]
Description=Quiz Application WebSocket Server
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html/quiz-app-websocket
ExecStart=/usr/bin/php /var/www/html/quiz-app-websocket/websocket-server.php

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/quiz-websocket/output.log
StandardError=append:/var/log/quiz-websocket/error.log
SyslogIdentifier=quiz-websocket

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/html/quiz-app-websocket/logs
ReadWritePaths=/var/log/quiz-websocket

# Resource limits
LimitNOFILE=10000
MemoryLimit=512M

[Install]
WantedBy=multi-user.target
